# climweb/docker/backup.cron
0 0 * * * cd /home/app/web && python3 manage.py dbbackup --clean --noinput && python3 manage.py mediabackup --clean --noinput  > /proc/1/fd/1 2>/proc/1/fd/2
