# climweb/docker/city_forecast.cron
0 * * * * cd /home/app/web && python3 manage.py generate_forecast > /proc/1/fd/1 2>/proc/1/fd/2
0 0 * * * cd /home/app/web && python3 manage.py clear_old_forecasts > /proc/1/fd/1 2>/proc/1/fd/2
