# climweb/docker/process_tasks.cron
*/15 * * * * source /climweb/venv/bin/activate && python3 /climweb/web/src/climweb/manage.py process_tasks --duration 900 > /proc/1/fd/1 2>/proc/1/fd/2
